# files/etc/services.d/php-fpm/run
#!/usr/bin/execlineb -P
s6-setuidgid nginx
php-fpm7
